# scripts/cacti/templates/s.m.a.r.t./smartd-snmp.cron
*/4 * * * root /etc/snmp/scripts/smartd-snmp_data.sh
